# deploy/systemd/ozybase.service
[Unit]
Description=OzyBase BaaS - Backend as a Service
Documentation=https://github.com/Xangel0s/OzyBase
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=OzyBase
Group=OzyBase

# Binary and working directory
ExecStart=/opt/OzyBase/ozybase
WorkingDirectory=/opt/OzyBase

# Auto restart
Restart=always
RestartSec=5

# Environment variables
EnvironmentFile=/opt/OzyBase/.env

# Resource limits (Essential for Realtime/SSE performance)
LimitNOFILE=65536
LimitNPROC=65536

# Security Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/OzyBase/data

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

